# deploy/auditd-rules.conf
# OpsGuard auditd reference rules.
# Purpose: forensic evidence only. Prevention is handled by least-privilege users,
# sudoers/ACL/read-only mounts, and application guard decisions.
#
# Install example:
#   sudo install -m 0640 deploy/auditd-rules.conf /etc/audit/rules.d/opsguard.rules
#   sudo augenrules --load

-w /etc/passwd -p wa -k opsguard_critical_identity
-w /etc/shadow -p wa -k opsguard_critical_identity
-w /etc/sudoers -p wa -k opsguard_critical_privilege
-w /etc/sudoers.d/ -p wa -k opsguard_critical_privilege
-w /etc/systemd/system/ -p wa -k opsguard_service_config
-w /var/lib/mysql/ -p wa -k opsguard_database_data
-w /var/lib/postgresql/ -p wa -k opsguard_database_data

# Destructive syscall watch for supported Linux auditd deployments.
-a always,exit -F arch=b64 -S unlink,unlinkat,rename,renameat,rmdir,chmod,fchmod,chown,fchown -k opsguard_mutating_syscall
-a always,exit -F arch=b32 -S unlink,unlinkat,rename,renameat,rmdir,chmod,fchmod,chown,fchown -k opsguard_mutating_syscall
